# TunisieTelecomCTF2K23/lemonldap-ng/fastcgi-server/systemd/llng-fastcgi-server.service
[Unit]
Description=FastCGI server for Lemonldap::NG websso system
After=network.target

[Service]
Type=forking
EnvironmentFile=/etc/default/llng-fastcgi-server
PIDFile=__FASTCGISOCKDIR__/llng-fastcgi-server.pid
ExecStart=__SBINDIR__/llng-fastcgi-server
ExecStop=-/sbin/start-stop-daemon --quiet --stop --retry QUIT/5 --pidfile ${PID}
KillMode=mixed

[Install]
Alias=llng-fastcgi-server.service
WantedBy=multi-user.target

