# TunisieTelecomCTF2K23/lemonldap-ng/debian/liblemonldap-ng-portal-perl.cron.d
#
# Regular cron jobs for LemonLDAP::NG Portal
#
7 *	* * *	www-data	[ -x /usr/share/lemonldap-ng/bin/purgeCentralCache ] && /usr/share/lemonldap-ng/bin/purgeCentralCache
